# justix-profe/debian/postinst
#!/bin/sh
# postinst script for justix-profe
#
# see: dh_installdeb(1)

set -e

creaProfe()
{
    useradd -m ${1} -s /bin/bash
    echo "${1}\n${1}" | passwd ${1} 2>/dev/null
}

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

# Diverted files
ORIG_FILE="/etc/default/arpon"
MY_FILE="/etc/default/arpon.justix"
MY_PACKAGE="justix-profe"

case "$1" in
    configure)


	# Afegim les diversions de fitxers
    	if dpkg-divert --package ${MY_PACKAGE} --add  --rename  --divert ${ORIG_FILE}.real ${ORIG_FILE} ; then
            ln -fs ${MY_FILE} ${ORIG_FILE}
        else
            echo "Unabe to divert file ${ORIG_FILE}"
        fi


	# Afegint l'usuari profe, fol i angles
	echo "\e[1m\e[93m [Justix] \e[95m Adding user profe \e[0m";
	`grep -q profe /etc/passwd` || creaProfe profe

	echo "\e[1m\e[93m [Justix] \e[95m Adding user angles \e[0m";
        `grep -q angles /etc/passwd` || creaProfe angles

	echo "\e[1m\e[93m [Justix] \e[95m Adding user fol \e[0m";
        `grep -q fol /etc/passwd` || creaProfe fol


        # Activant arpon
        systemctl enable arpon.service
        systemctl start arpon.service
	    

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
